# rfchain-hud/deploy/rfchain-hud.service
[Unit]
Description=RFChain HUD - Jarvis Signal Analyzer
Documentation=https://github.com/cvalentine99/RFChain
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=rfchain
Group=rfchain
WorkingDirectory=/opt/rfchain-hud

# Environment setup
Environment=NODE_ENV=production
Environment=PORT=3007
Environment=RFCHAIN_INSTALL_DIR=/opt/rfchain-hud
Environment=RFCHAIN_CONDA_ENV=rfchain-hud
Environment=RFCHAIN_DB_TIMEOUT=10
EnvironmentFile=/opt/rfchain-hud/.env.local

# Pre-start dependency checks
# -: prefix means don't fail if script returns non-zero (warnings only)
# Without prefix, critical failures will prevent service start
ExecStartPre=/opt/rfchain-hud/deploy/check-deps.sh

# Conda environment activation and server start
ExecStart=/bin/bash -c 'source /home/rfchain/miniconda3/etc/profile.d/conda.sh && conda activate rfchain-hud && exec node dist/server/index.js'

# Graceful shutdown
ExecStop=/bin/kill -SIGTERM $MAINPID
TimeoutStopSec=30

# Restart policy
Restart=on-failure
RestartSec=10
StartLimitIntervalSec=60
StartLimitBurst=3

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ReadWritePaths=/opt/rfchain-hud/uploads /opt/rfchain-hud/analysis_output /opt/rfchain-hud/logs

# Resource limits
LimitNOFILE=65536
MemoryMax=4G
CPUQuota=400%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rfchain-hud

[Install]
WantedBy=multi-user.target
