# rfchain-hud/deploy/rfchain-hud.service
[Unit]
Description=RFChain HUD - Jarvis Signal Analyzer
Documentation=https://github.com/cvalentine99/RFChain
After=network.target network-online.target mysql.service
Wants=network-online.target

[Service]
Type=simple
User=cvalentine
Group=cvalentine
WorkingDirectory=/home/cvalentine/RFChain/rfchain-hud

# Environment setup
Environment=NODE_ENV=production
Environment=PORT=3007
EnvironmentFile=/home/cvalentine/RFChain/rfchain-hud/.env.local

# Pre-start dependency checks (optional - remove if check-deps.sh not installed)
# ExecStartPre=/home/cvalentine/RFChain/rfchain-hud/deploy/check-deps.sh

# Start the production server (dist/index.js is the bundled output)
ExecStart=/usr/bin/node /home/cvalentine/RFChain/rfchain-hud/dist/index.js

# Graceful shutdown
ExecStop=/bin/kill -SIGTERM $MAINPID
TimeoutStopSec=30

# Restart policy
Restart=on-failure
RestartSec=10
StartLimitIntervalSec=60
StartLimitBurst=3

# Resource limits
LimitNOFILE=65536

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rfchain-hud

[Install]
WantedBy=multi-user.target
